# scripts/systemd/nostr-codex-runner.service
[Unit]
Description=Nostr Codex Runner
After=network.target

[Service]
Type=simple
User=%i
Environment=NCR_CONFIG=%h/.config/nostr-codex-runner/config.yaml
ExecStart=%h/.local/bin/nostr-codex-runner run -config ${NCR_CONFIG}
WorkingDirectory=%h
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
